# etc/fim-pg_repack.cron
*/30 * * * * postgres /usr/pgsql-9.4/bin/pg_repack fim 
